# packaging/fifine-deck.service
[Unit]
Description=fifine Control Deck (Linux) — key daemon
After=graphical-session.target

[Service]
# Adjust the path if you cloned elsewhere.
ExecStart=/usr/bin/python3 -m fifine_deck --headless
WorkingDirectory=%h/FifineControlDeck
Restart=on-failure
RestartSec=2
# Give wpctl / playerctl / xdg-open the right session bus + display.
Environment=XDG_RUNTIME_DIR=%t

[Install]
WantedBy=default.target
